# build/rehl8/go-http-server.service
[Unit]
Description=go-http server

[Service]
Type=simple
ExecStart=/usr/bin/go-http-server -config /etc/go-http-server/config.json

[Install]
WantedBy=multi-user.target